# files/supervisor.conf


[program:gunicorn]
command=/usr/local/bin/gunicorn --bind 0.0.0.0:8192 --reload rsum.wsgi 
directory=/srv/rsum/rsum
user=root
autostart=true
autorestart=true
redirect_stderr=true
